# i-am-making-a-bad-decision.service
[Unit]
Description=A website blocker
After=multi-user.target graphical.target network-online.target
[Service]
Type=simple
Restart=always
RestartSec=10
TimeoutSec=2
WorkingDirectory=/home/charlie/src/python/hb
ExecStart=/home/charlie/src/python/hb/main.py
User=root
[Install]
WantedBy=multi-user.target